# coned-scraper/rootfs/etc/s6-overlay/s6-rc.d/init-coned/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ConEd Scraper - Init script
# ==============================================================================

bashio::log.info "ConEd Scraper initializing..."
