# coned_scraper/rootfs/etc/s6-overlay/s6-rc.d/web/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# ConEd Scraper - Next.js frontend (ingress-aware via assetPrefix injection)
# ==============================================================================

bashio::log.info "Starting ConEd Scraper Web UI..."

export NODE_ENV=production
export PORT=3000
export HOSTNAME=0.0.0.0
export API_BASE_URL=http://127.0.0.1:8000

# Inject ingress path into Next.js built files so assets load correctly
INGRESS_PATH="$(bashio::addon.ingress_entry 2>/dev/null || true)"
export INGRESS_PATH
cd /app/web
if [ -n "$INGRESS_PATH" ]; then
  bashio::log.info "Setting assetPrefix to: $INGRESS_PATH"
  node inject-ingress.js
fi
exec node server.js
